# drive-monitor.service
[Unit]
Description=Drive usage and health monitor daemon

[Service]
ExecStart=/usr/local/sbin/drive-monitor
